# roles/containers/templates/https.socket.j2
# {{ ansible_managed }}
# {{ ansible_fqdn }}
# {{{{ template_path|basename + lookup('pipe', 'git log --format=",%%h %%ad %%ae" -1 --date=format:"%%Y/%%m/%%d %%H:%%M" ' + template_fullpath|quote)|default(",UNCOMMITED", True) }}}}

[Socket]
# forward 443 to 8443 using firewalld
ListenStream=0.0.0.0:8443
FileDescriptorName=websecure
Service=traefik.service

[Install]
WantedBy=sockets.target
